# magang_infer/Maintenance/supervisor/supervisord.conf
[unix_http_server]
# todo 修改supervisord配置文件
file=/tmp/supervisor.sock   ; the path to the socket file
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
username=admin              ; default is no username (open server)
password=123               ; default is no password (open server)

[inet_http_server]         ; inet (TCP) server disabled by default
port=192.168.100.3:9001       ; ip_address:port specifier, *:port for all iface
username=admin              ; default is no username (open server)
password=123               ; default is no password (open server)

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
serverurl=http://192.168.100.3:9001 ; use an http:// url to specify an inet socket
username=admin              ; should be same as in [*_http_server] if set
password=123                ; should be same as in [*_http_server] if set


[program:new_infer_runner]  ;程序名称
user=hongtai  ;执行程序的用户
command=/usr/bin/python3 /home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/MyAlgorithmRunner.py  ;执行的命令
directory=/home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/ ;命令执行的目录
stopsignal=TERM  ;重启时发送的信号
autostart=true  
autorestart=true  ;是否自动重启
stdout_logfile=/home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/Log/run_logs/stdout.log  ;标准输出日志位置
stdout_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
stderr_logfile=/home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/Log/run_logs/stderr.log  ;标准错误日志位置
stderr_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)

[program:image_delete]  ;程序名称
user=hongtai  ;执行程序的用户
command=/usr/bin/python3 /home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/UtilObject/ImageDeleteUtil.py  ;执行的命令
directory=/home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/ ;命令执行的目录
stopsignal=TERM  ;重启时发送的信号
autostart=true  
autorestart=true  ;是否自动重启
stdout_logfile=/home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/Log/del_logs/del_out.log  ;标准输出日志位置
stdout_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
stderr_logfile=/home/hongtai/yolo/magang_infer/SteelDefectDetection-magang/magang_infer/Log/del_logs/del_err.log  ;标准错误日志位置
stderr_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)

# [program:image_backup]  ;程序名称
# user=hongtai  ;执行程序的用户
# command=/home/hongtai/anaconda3/bin/python3 /home/hongtai/workspace/general_code/image_backup.py  ;执行的命令
# directory=/home/hongtai/workspace/general_code/ ;命令执行的目录
# stopsignal=TERM  ;重启时发送的信号
# autostart=true  
# autorestart=true  ;是否自动重启
# stdout_logfile=/home/hongtai/workspace/general_code/backup_log/bkup_out.log  ;标准输出日志位置
# stdout_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
# stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
# stderr_logfile=/home/hongtai/workspace/general_code/backup_log/bkup_err.log  ;标准错误日志位置
# stderr_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
# stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)


